# deploy/homehub.service
[Unit]
Description=HomeHub
After=network.target sound.target

[Service]
User=pi
WorkingDirectory=/home/pi/homehub
Environment=PYTHONPATH=/home/pi/homehub
ExecStart=/home/pi/homehub/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always

[Install]
WantedBy=multi-user.target
